# deploy/vanna.service
[Unit]
Description=Vanna API Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/vanna_v2
Environment="PATH=/home/ubuntu/vanna_v2/venv/bin"
EnvironmentFile=/home/ubuntu/vanna_v2/.env
ExecStart=/home/ubuntu/vanna_v2/venv/bin/gunicorn -w 4 -k uvicorn.workers.UvicornWorker app:server.app --bind 0.0.0.0:8080
Restart=always

[Install]
WantedBy=multi-user.target